# ansible/roles/monitoring_setup/templates/backup.timer.j2
[Unit]
Description=Run TG-RSS Backup daily at 2 AM
Requires=tg-rss-backup.service

[Timer]
OnCalendar=daily
Persistent=true
Unit=tg-rss-backup.service

[Install]
WantedBy=timers.target
