# ansible/ansible/playbooks/roles/monitoring_setup/templates/app-monitor.timer.j2
[Unit]
Description=Run TG-RSS Application Monitor every 5 minutes
Requires=app-monitor.service

[Timer]
OnBootSec=5min
OnUnitActiveSec=5min
Unit=app-monitor.service

[Install]
WantedBy=timers.target
