# ansible/ansible/playbooks/roles/monitoring_setup/templates/app-monitor.service.j2
[Unit]
Description=TG-RSS Application Monitor
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/tg-rss-monitor.sh
User=root
