# ansible/ansible/playbooks/roles/monitoring_setup/templates/backup.service.j2
[Unit]
Description=TG-RSS Database Backup
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/tg-rss-backup.sh
User=root
